# infra/systemd/ptt-fleet-api.service
[Unit]
Description=PTT Fleet API Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ptt-fleet
Group=ptt-fleet
WorkingDirectory=/opt/ptt-fleet/api-server
EnvironmentFile=/etc/ptt-fleet/api-server.env
ExecStart=/opt/ptt-fleet/api-server/api-server
Restart=on-failure
RestartSec=5s
TimeoutStartSec=30s
TimeoutStopSec=20s
KillSignal=SIGTERM

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ptt-fleet-api

[Install]
WantedBy=multi-user.target
